# e2e/test-log-generator.service
[Unit]
Description=Test Log Generator for E2E Testing
After=systemd-journald.service

[Service]
Type=simple
ExecStart=/app/e2e/test-log-generator.sh
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target

